# deploy/systemd/a-share-quant-eod.service
[Unit]
Description=A-share swing quant weekday EOD observation
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/curl --config /etc/a-share-quant/eod.curl.conf
Restart=on-failure
RestartSec=300
