# deploy/systemd/a-share-quant-backup.service
[Unit]
Description=Back up A-share quant SQLite data outside the Docker volume
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/a-share-swing-quant
ExecStart=/usr/bin/docker compose exec -T api python -m quant_system.backup --source /app/data/quant_system.db --output /app/backups --retention-days 30 --artifact /app/data/research
